# examples/deploy/vps/ray-gateway.service
[Unit]
Description=Ray Gateway
After=network.target

[Service]
Type=simple
User=ray
WorkingDirectory=/srv/ray
EnvironmentFile=/etc/ray/ray.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /srv/ray/apps/gateway/dist/index.js --config /etc/ray/ray.json
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
LockPersonality=true
RestrictSUIDSGID=true
UMask=027
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
