# examples/deploy/vps/ray-gateway.service
[Unit]
Description=Ray Gateway
Wants=ray-llama-cpp.service
After=network.target ray-llama-cpp.service
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=ray
WorkingDirectory=/srv/ray
EnvironmentFile=/etc/ray/ray.env
Environment=NODE_ENV=production
ExecStart=/usr/local/bin/bun /srv/ray/apps/gateway/dist/index.js --config /etc/ray/ray.json
Restart=always
RestartSec=2
LogRateLimitIntervalSec=30s
LogRateLimitBurst=200
TimeoutStopSec=35
KillSignal=SIGTERM
KillMode=mixed
OOMPolicy=stop
OOMScoreAdjust=-250
TasksMax=128
CPUAccounting=true
CPUWeight=200

MemoryAccounting=true
IOAccounting=true
MemoryHigh=640M
MemoryMax=896M
MemorySwapMax=128M

StateDirectory=ray
StateDirectoryMode=0750
NoNewPrivileges=true
CapabilityBoundingSet=
SystemCallArchitectures=native
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectClock=true
ProtectHostname=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
UMask=077
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
